# rootfs/etc/services.d/live555proxyserver/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# SmartHouse Hass.io Add-ons: Live555 RTSP Proxy Server
# Runs the Live555 RTSP Proxy Server
# ==============================================================================

LOG_LEVEL=$(bashio::config 'log_level')
CONFIG_PATH=/data/options.json
STREAMS=$(jq --raw-output ".streams[]" $CONFIG_PATH)

export LOG_LEVEL

bashio::log.info "Starting Live555 RTSP Proxy Server..."
exec live555ProxyServer ${STREAMS[*]}
